# proxmox-gui.service
[Unit]
Description=Proxmox Lab GUI
After=network.target

[Service]
Type=notify
User=root
Group=root
WorkingDirectory=/root/rdp-gen
Environment="PATH=/root/rdp-gen/venv/bin"
EnvironmentFile=-/root/rdp-gen/.env
ExecStart=/root/rdp-gen/venv/bin/gunicorn --config /root/rdp-gen/gunicorn.conf.py --chdir /root/rdp-gen app:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3

# Security options
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
